# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:one-api]
command=/one-api
priority=10
autostart=true
autorestart=unexpected
startsecs=5
startretries=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:picoclaw]
command=/usr/local/bin/picoclaw-gateway.sh
priority=20
autostart=true
autorestart=unexpected
startsecs=5
startretries=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
