# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:nginx]
command=nginx -g "daemon off;"
priority=5
autostart=true
autorestart=unexpected
startsecs=3
startretries=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:one-api]
command=/bin/sh -c "PORT=3001 exec /one-api"
priority=10
autostart=true
autorestart=unexpected
startsecs=5
startretries=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:picoclaw-launcher]
command=/usr/local/bin/picoclaw-launcher -console -public -no-browser
priority=20
autostart=true
autorestart=unexpected
startsecs=5
startretries=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:test-runner]
command=/usr/local/bin/python3 /app/skills/test_runner/scripts/test_runner.py
priority=25
priority=20
autostart=true
autorestart=unexpected
startsecs=5
startretries=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
